# scripts/com.inventory.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.inventory.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string><%= serverPath %></string>
    </array>
    <key>WorkingDirectory</key>
    <string><%= appDir %></string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string><%= appDir %>/server.log</string>
    <key>StandardErrorPath</key>
    <string><%= appDir %>/server.log</string>
</dict>
</plist>
